# build/production/schedule.cron
*/30 * * * * php /var/www/current/protected/t4.php emailSender >> /var/log/cron.log &
